# FOOD_APP/myapp/deploy/inj-reviews-frontend.service
[Unit]
Description=INJ Reviews frontend static proxy
After=network-online.target inj-reviews-backend.service
Wants=network-online.target
Requires=inj-reviews-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/inj-reviews
EnvironmentFile=/home/ubuntu/inj-reviews/.env
Environment=BASE_PATH=/reviews
Environment=HOST=127.0.0.1
Environment=PORT=8081
Environment=FRONTEND_PORT=8081
Environment=DIST_DIR=/home/ubuntu/inj-reviews/dist
Environment=BACKEND_HOST=127.0.0.1
Environment=BACKEND_PORT=8787
ExecStart=/usr/bin/python3 /home/ubuntu/inj-reviews/deploy/static_proxy_server.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
